# root/etc/services.d/pyload/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc pyload --userdir /config --storagedir /downloads
